# tools/monitor/com.snorkel.terminus-imessage-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.snorkel.terminus-imessage-control</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-W</string>
        <string>__LISTENER_APP__</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/.review-scratch/monitor/imessage-listener-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/.review-scratch/monitor/imessage-listener-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>REPO_ROOT</key>
        <string>__REPO_ROOT__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>
